# mesos.mesos-slave.upstart
# mesos slave upstart
description "Mesos Slave"

env DEFAULTFILE=/etc/default/mesos

start on runlevel [2345]
respawn

pre-start script
  RUN_MESOS_SLAVE=no
  if [ -f "$DEFAULTFILE" ]; then
    . "$DEFAULTFILE"
  fi
  if [ "x$RUN_MESOS_MASTER" != xyes ]; then
    stop
    exit 0
  fi
  ulimit -n 8192
end script

script
  exec /usr/sbin/mesos-slave --conf=/etc/mesos/conf
end script


