# ubuntu/master.upstart
description "mesos master"

start on runlevel [2345]
respawn

script
  [ ! -f /etc/default/mesos ]        || . /etc/default/mesos
  [ ! -f /etc/default/mesos-master ] || . /etc/default/mesos-master
  [ -z "${ULIMIT:-}" ] || ulimit $ULIMIT
  exec /usr/local/sbin/mesos-master \
    --zk=file:///etc/mesos/zk --port="$MASTER_PORT" --cluster="${CLUSTER:-}" \
    >> "$LOGS"/master.log 2>&1
end script

